# app-misc/sabayon-music/files/music
#!/sbin/runscript
# Copyright 2006 SabayonLinux
# Distributed under the terms of the GNU General Public License v2

depend() {
	before alsasound
}


start() {

	cmdline_music_exist=$(cat /proc/cmdline | grep "music")
	cmdline_mute_exist=$(cat /proc/cmdline | grep "sound=mute")
	sound_device_available=$(lspci | grep -i "audio")

	if [ -n "$sound_device_available" ]; then

		# Setting up volumes
                /usr/bin/amixer -q set "Analog Front" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "Analog Rear" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "Analog Side" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "Analog Center/LFE" 70% unmute &> /dev/null
                /usr/bin/amixer -q set Master 70% unmute &> /dev/null
                /usr/bin/amixer -q set PCM 70% unmute &> /dev/null
                /usr/bin/amixer -q set Front 70% unmute &> /dev/null
                /usr/bin/amixer -q set Surround 70% unmute &> /dev/null
                /usr/bin/amixer -q set Center 70% unmute &> /dev/null
                /usr/bin/amixer -q set LFE 70% unmute &> /dev/null
                /usr/bin/amixer -q set "Analog Mix" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "Aux" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "Aux2" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "PCM Center" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "PCM Front" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "PCM LFE" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "PCM Side" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "PCM Surround" 70% unmute &> /dev/null
                /usr/bin/amixer -q set "Side" 70% unmute &> /dev/null
                /usr/bin/amixer -q set 'Input Gain' 50% &> /dev/null
                /usr/bin/amixer -q set FM 60% unmute &> /dev/null
                alsactl store &> /dev/null
		touch /etc/amixer_has_been_run
	fi



	if [ -n "$cmdline_music_exist" ] && [ -z "$cmdline_mute_exist" ] && [ -n "$sound_device_available" ]; then
		ebegin "Starting music"
			# get memory
			memory_avail=$(echo $(free -m | grep "Mem") | cut -d" " -f 2)
			if [ "$memory_avail" -gt "384" ]; then
				# copy the image in ram
				cp /usr/share/sounds/boot.ogg /var/tmp/music.ogg &> /dev/null
				nice -n -15 ogg123 -p 30 -q /var/tmp/music.ogg &> /dev/null &
			else
				nice -n -15 ogg123 -p 30 -q /usr/share/sounds/boot.ogg &> /dev/null &
			fi
		eend 0
	else
		ebegin "Ok... music disabled"
		eend 0
	fi
}
