# sys-fs/mdadm/files/mdraid.rc-3.1.1
#!/sbin/runscript
# Copyright 1999-2007 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/sys-fs/mdadm/files/mdraid.rc-3.1.1,v 1.1 2010/01/24 22:20:57 robbat2 Exp $

depend() {
	if [ ! -e /sbin/functions.sh ] ; then
		# on baselayout-1 this could cause
		# dependency-cycles with checkroot (before *)
		before checkfs fsck
		after modules
	fi
}

run_addon() {
	local f=/lib/rcscripts/addons/$1.sh
	if [ -r "$f" ]; then
		( . "$f" )
	fi
}
start_addon() {
	run_addon "$1"-start
}
stop_addon() {
	run_addon "$1"-stop
}

start() {
	if [ -e /sbin/functions.sh ] ; then
		eerror "The ${SVCNAME} init script is written for baselayout-2"
		eerror "Please do not use it with baselayout-1"
		return 1
	fi

	start_addon raid
}

stop() {
	stop_addon raid
}
