# app-emulation/virtualbox-bin/files/vboxwebsrv-initd
#!/sbin/runscript
# Copyright 1999-2005 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/app-emulation/virtualbox-bin/files/vboxwebsrv-initd,v 1.4 2009/03/21 10:42:14 patrick Exp $

depend() {
	need net
}

start() {
	#checkconfig || return $?

	ebegin "Starting VirtualBox WebService"
	start-stop-daemon --start --quiet --pidfile /var/run/vboxwebsrv.pid \
		--background --make-pidfile --exec /opt/bin/vboxwebsrv \
		-- -H ${VBOXWEBSRV_HOST} -p ${VBOXWEBSRV_PORT} -t ${VBOXWEBSRV_TIMEOUT} -i ${VBOXWEBSRV_INTERVAL}
	eend $?
}

stop() {
        ebegin "Stopping VirtualBox WebService"
        start-stop-daemon --stop --quiet --pidfile /var/run/vboxwebsrv.pid
        eend $?
}
