# app-misc/sabayon-mce/files/1.0/init.d/sabayon-mce
#!/sbin/runscript
# Copyright 2009 Sabayon Linux

depend() {
	after mtab
	before xdm
	# In this way autologin settings are overwritten
	after sabayonlive
}

start() {

	source /sbin/sabayon-functions.sh
	sabayonmce_avail=$(cat /proc/cmdline | grep sabayonmce)

	if [ -n "$sabayonmce_avail" ]; then
		ebegin "Sabayon Media Center mode enabled"

		# setup default .dmrc
		echo "[Desktop]" > /var/sabayonmce/.dmrc
		echo "Session=sabayon-mce" >> /var/sabayonmce/.dmrc
		chown sabayonmce /var/sabayonmce/.dmrc

		SABAYON_USER="sabayonmce"
		sabayon_setup_autologin
	else
		ebegin "Sabayon Media Center mode disabled"
		sabayon_disable_autologin
	fi

	eend 0

}

