# net-misc/networkmanager/files/NetworkManager-setup
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need localmount
	before NetworkManager
}

start() {
	if get_bootparam "nonm" ; then
		touch /etc/.nonm
	fi
}
