# net-misc/networkmanager/files/NetworkManager-setup
#!/sbin/runscript
# Copyright 1999-2011 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: $

depend() {
	need localmount
	before NetworkManager
}

start() {
	if get_bootparam "nonm" ; then
		local nmfile=/etc/NetworkManager/NetworkManager.conf
		ebegin "Setting managed=false in ${nmfile}"
		sed -i "s:^managed=.*:managed=false:" "${nmfile}" || eend 1
		eend 0
	fi
}
