# x11-apps/xinit/files/xdm-setup.initd-1
#!/sbin/runscript
# Copyright 1999-2009 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/x11-apps/xinit/files/xdm-setup.initd-1,v 1.6 2009/10/06 16:44:38 williamh Exp $

depend() {
	need localmount
}

start() {
	if get_bootparam "nox" ; then
		touch /etc/.noxdm
	fi
}
