# app-misc/sabayonlive-tools/files/2.2/sabayonlive
#!/sbin/runscript
# Copyright 2004-2009 Fabio Erculiani
# Distributed under the terms of the GNU General Public License v2 

. /sbin/sabayon-functions.sh

depend() {
	after x-setup
	before hald
	before xdm
}

prepare_fluxbox_locked_autologin() {

	# change default wm to fluxbox
	echo "[Desktop]" > /home/${LIVE_USER}/.dmrc
	echo "Session=fluxbox" >> /home/${LIVE_USER}/.dmrc
	chown ${LIVE_USER} /home/${LIVE_USER}/.dmrc
	# Change passwords for security
	echo root:\`pwgen -s 16\` | chpasswd  > /dev/null 2>&1
	echo ${LIVE_USER}:\`pwgen -s 16\` | chpasswd  > /dev/null 2>&1

}

my_setup_desktop() {

	# create LIVE_USER if it does not exist
	sabayon_setup_live_user "${LIVE_USER}" "1000"
	if [ "${?}" = "1" ]; then
		# if user is already available, then setup skel
		# Copy ${LIVE_USER} directory
		rm -rf /home/${LIVE_USER}
		cp /etc/skel /home/${LIVE_USER} -Rp
		chown ${LIVE_USER}:users /home/${LIVE_USER} -R
	fi

	liveinst_desktop="/usr/share/applications/liveinst.desktop"
	liveinst_desktop_name="$(basename ${liveinst_desktop})"
	if [ -f "${liveinst_desktop}" ]; then
		[[ -d "/home/${LIVE_USER}/Desktop" ]] || mkdir -p /home/${LIVE_USER}/Desktop 
		cp "${liveinst_desktop}" /home/${LIVE_USER}/Desktop
		chown ${LIVE_USER}:users /home/${LIVE_USER}/Desktop -R
		chmod +x /home/${LIVE_USER}/Desktop/${liveinst_desktop_name}
		rm -f /etc/skel/Desktop/Anaconda*.desktop
		rm -f /home/${LIVE_USER}/Desktop/Anaconda*.desktop
	fi

	# Disable memory eating services
	rm -f /etc/xdg/autostart/hplip-systray.desktop
	rm -f /etc/xdg/autostart/beagle-search-autostart.desktop
	rm -f /etc/xdg/autostart/tracker*.desktop
	rm -f /etc/xdg/autostart/magneto.desktop
	rm -f /etc/xdg/autostart/beagled-autostart.desktop
	rm -f /usr/share/autostart/magneto.desktop
	rm -f /usr/share/autostart/nepomukserver.desktop

        # Remove broken entries in /etc/mtab
        sed -i '/.*newroot.*/d' /etc/mtab

	# Create /media for removable devices
	if [ ! -d /media ]; then
		mkdir /media
		chmod 755 /media
	fi

	# Add sudo to gparted exec=
	gparted_file="/home/${LIVE_USER}/Desktop/gparted.desktop"
	if [ -f "${gparted_file}" ]; then
		sed -i 's/Exec=/Exec=sudo /' "${gparted_file}"
	fi

	# create /overlay, this way df -h won't bitch
	[[ -d "/overlay" ]] || mkdir /overlay

}

my_setup_password() {

	cmdline_autoscramble_exist=$(cat /proc/cmdline | grep autoscramble)
	if [ -n "$cmdline_autoscramble_exist" ]; then
		echo
		echo -e "\E[33;36m * \E[0m\E[01;36m Autoscrambling root passwords for   S E C U R I T Y"
		echo root:\`pwgen -s 16\` | chpasswd  > /dev/null 2>&1
		echo ${LIVE_USER}:\`pwgen -s 16\` | chpasswd  > /dev/null 2>&1
		echo
	fi

}

my_setup_session() {

	cmdline_session_exist=$(cat /proc/cmdline | grep session=)
	if [ -n "$cmdline_session_exist" ]; then
		# Session exists, should be: kde,gnome,fluxbox
                for word in `cat /proc/cmdline` ; do
			case $word in
				session=*)
					my_sess=$(echo $word | cut -d "=" -f 2)
					# change default DE to KDE
					echo "[Desktop]" > /home/${LIVE_USER}/.dmrc
					echo "Session="$my_sess >> /home/${LIVE_USER}/.dmrc
					chown ${LIVE_USER} /home/${LIVE_USER}/.dmrc
				;;
			esac
		done
	fi

}

my_setup_xsistence() {

	persistent_exists=$(cat /proc/cmdline | grep "xsistence")
	if [ -z "$persistent_exists" ]; then
		return 0
	fi

	xsistence_mb=$(cat /proc/cmdline | grep "xsistence_mb=")
	if [ -n "$xsistence_mb" ];then
		for word in `cat /proc/cmdline` ; do
		case $word in
			xsistence_mb=*)
			xsistence_img_mb=$(echo $word | cut -d "=" -f 2)
			;;
		esac
		done
	else
		xsistence_img_mb=28
	fi

	# Persistent Mode enabled
	echo -e "\E[33;36m * \E[0m \E[01;32m Enabling XsistenCe Mode (logfile: /tmp/xsistence.log)... \E[0m"
	echo -en "\E[33;36m * \E[0m \E[01;32m Please plug-in your USB Memory Stick now (10 seconds)... \E[0m"
	sleep 10
	echo -en "Waiting for devices..."
	sleep 6
	echo -e "Scanning."
	sleep 1

	do_break="0"
	xsistence_found="0"

	for file in /sys/block/*/device; do

		if [ "$do_break" = "1" ]; then
			break
		fi

		is_floppy=$(echo $file | grep "fd[0-99]")
		is_sr=$(echo $file | grep "sr[0-99]")
		if [ -n "$is_floppy" ] || [ -n "$is_sr" ]; then
			continue
		fi

		# if it's not a floppy or an external recorder, scrub it
		device=$(echo $file | cut -d/ -f4)
		is_removable="0"
		if [ -e /sys/block/$device/removable ]; then
			is_removable=$(cat /sys/block/$device/removable)
		fi
		if [ "$is_removable" != "1" ]; then
			continue
		fi

		# Counting partitions...
		partitions_list=$(ls /dev/$device*)

		for word in $partitions_list; do

			# Mount them and verify the correct operation
			mkdir -p /mnt/xsistence &> /dev/null
			mount $word /mnt/xsistence &> /dev/null

			if [ "$?" != "0" ]; then
				echo "Tried to mount "$word" without luck" >> /tmp/xsistence.log
				continue
			fi
			echo "Mounted "$word" for test" >> /tmp/xsistence.log

			# Look for "xsistence-mode" file (case insensitive)
			file_found=$(find /mnt/xsistence -maxdepth 1 -nowarn -iname xsistence-mode)
			if [ -z "$file_found" ];then
				echo "xsistence-mode not found on "$word >> /tmp/xsistence.log
				umount /mnt/xsistence -f &> /dev/null
				rmdir /mnt/xsistence /&> /dev/null
				continue
			fi
			
			# test r/w capabilities
			touch /mnt/xsistence/.touch-test &> /dev/null
			if [ "$?" != "0" ] || [ ! -f "/mnt/xsistence/.touch-test" ]; then
				echo "No read-write capabilities for " $word >> /tmp/xsistence.log
				umount /mnt/xsistence -f &> /dev/null
				rmdir /mnt/xsistence /&> /dev/null
				continue
			fi

			rm /mnt/xsistence/.touch-test
			# Look for existing image
			if [ -e /mnt/xsistence/sabayon-xsistence.ext3 ]; then

				echo "/mnt/xsistence/sabayon-xsistence.ext3 on "$word " exists" >> /tmp/xsistence.log
				
				# found ! mount it
				mount -o loop /mnt/xsistence/sabayon-xsistence.ext3 /home &> /dev/null
				if [ "$?" = "0" ]; then
					echo "/mnt/xsistence/sabayon-xsistence.ext3 mounted using "$word >> /tmp/xsistence.log
					echo  -e "\t \t Found XsistenCe on "$word". Mounted and fully working."
					do_break="1"
					xsistence_found="1"
					break
				else
					echo "/mnt/xsistence/sabayon-xsistence.ext3 cannot be mounted  using " $word >> /tmp/xsistence.log
					umount /mnt/xsistence -f &> /dev/null
					rmdir /mnt/xsistence /&> /dev/null
					echo  -e "\t \t Found XsistenCe on "$word". But the image does not work."
					continue
				fi
				
			fi

			echo "/mnt/xsistence/sabayon-xsistence.ext3 on "$word "does not exist" >> /tmp/xsistence.log
			# create disk image and test available space
			# need to find another way to get avaiable space
			
			# output
			echo  -en "\t \t Trying to create the image on "$word" ..."
			dd if=/dev/zero of=/mnt/xsistence/sabayon-xsistence.ext3 count=$xsistence_img_mb bs=1024000 &> /dev/null

			lofine="0"
			for lodev in "7 6 5 4 3 2"; do
				loopdevice=${lodev}
				losetup /dev/loop${lodev} /mnt/xsistence/sabayon-xsistence.ext3
				if [ "$?" = "0" ]; then
					lofine="1"
					break
				fi
			done

			if [ "$lofine" = "0" ]; then
				echo -e "Sorry, no free loopback device"
				echo "No free loopback device for "$word >> /tmp/xsistence.log
				continue
			fi

			# loopback ok, formatting...
			mkfs.ext3 /dev/loop$loopdevice &> /dev/null

			if [ "$?" != "0" ]; then
				echo -e "Sorry, unable to format" $xsistence_img_mb
				echo "Cannot format the image on " $word  >> /tmp/xsistence.log
				rm -rf /mnt/xsistence/sabayon-xsistence.ext3 &> /dev/null
				continue
			fi

			echo "Image formatted and ready to get filled " $word "using loopback" $loopdevice >> /tmp/xsistence.log
			mkdir /mnt/loopbacktemp
			mount -o loop /mnt/xsistence/sabayon-xsistence.ext3 /mnt/loopbacktemp

			if [ "$?" != "0" ]; then
				echo -e "Sorry, cannot mount the image"
				echo "Image not mounted" $word "using loopback" $loopdevice >> /tmp/xsistence.log
				continue
			fi

			echo "Image formatted and ready to get filled " $word "using loopback" $loopdevice >> /tmp/xsistence.log
			cp -Rap /home/* /mnt/loopbacktemp/
			mount --move /mnt/loopbacktemp/ /home/

			if [ "$?" != "0" ]; then
				echo -e "Sorry, cannot move mount point"
				echo "Cannot move mount point for" $word "and loopback" $loopdevice  >> /tmp/xsistence.log
				continue
			fi

			echo -e "image created and mounted !"
			echo "Image created and mounted using" $word "and loopback" $loopdevice >> /tmp/xsistence.log
			do_break="1"
			xsistence_found="1"
			break

		done
	
	done

	if [ "$xsistence_found" = "0" ]; then
		# Persistent Mode disabled
		echo -e "\E[33;36m * \E[0m \E[01;32m XsistenCe is not enabled or no devices found. \E[0m"
	fi

}

my_setup_keymap() {

	# Check if KEYMAP is forced by cmdline
	cmdline_keymap_exist=$(cat /proc/cmdline | grep -i "KEYMAP=")
	cmdline_keymap_isolinux_exist=$(cat /proc/cmdline | grep -i "console-setup/layoutcode=")

	if [ -n "$cmdline_keymap_isolinux_exist" ]; then
		cmdline_keymap_exist="console-setup"
	fi

        if [ -n "$cmdline_keymap_exist" ];then

		if [ "$cmdline_keymap_exist" == "console-setup" ]; then
			# detect keymap
			for word in `cat /proc/cmdline` ; do
				case $word in
					console-setup/layoutcode=*)
					keymap_toset=$(echo $word | cut -d "=" -f 2)
					;;
					console-setup/modelcode=*)
					keymap_toset_model="-$(echo $word | cut -d "=" -f 2)"
					;;
				esac
			done
		else

			# detect keymap
			for word in `cat /proc/cmdline` ; do
				case $word in
					KEYMAP=*)
					keymap_toset=$(echo $word | cut -d "=" -f 2)
					;;
					keymap=*)
					keymap_toset=$(echo $word | cut -d "=" -f 2)
					;;
				esac
			done

		fi

		if [ -n "$keymap_toset" ]; then
			/sbin/keyboard-setup $keymap_toset$keymap_toset_model &> /dev/null
		fi
	
	fi

}

my_wait_opengl_config() {

	echo -en "\E[33;36m * \E[0m \E[01;32m Waiting for Graphical Configurator to finish... \E[0m"
	while [ -e "/etc/x-setup-configuration-running" ]
           do
	      # waiting...
	      sleep 1
	   done
	echo -e "Done"

}

my_setup_locale() {

        cmdline_lang_exist=$(cat /proc/cmdline | grep -i "lang=")
        cmdline_locale_exist=$(cat /proc/cmdline | grep -i "locale=")

	if [ -n "$cmdline_locale_exist" ]; then
		cmdline_lang_exist="locale-setup"
	fi

        if [ -n "$cmdline_lang_exist" ];then

		if [ "$cmdline_lang_exist" == "locale-setup" ]; then
			# detect keymap
			for word in `cat /proc/cmdline` ; do
				case $word in
					locale=*)
					lang_toset=$(echo $word | cut -d "=" -f 2)
					;;
				esac
			done
		else
			for word in `cat /proc/cmdline` ; do
				case $word in
					LANG=*)
					lang_toset=$(echo $word | cut -d "=" -f 2)
					;;
					lang=*)
					lang_toset=$(echo $word | cut -d "=" -f 2)
					;;
				esac
			done
		fi
		
		# Setup Language
		if [ -n "$lang_toset" ]; then
			/sbin/language-setup $lang_toset &> /dev/null
		fi

        fi

}

my_setup_xsistence_stop() {

	persistent_exists=$(cat /proc/cmdline | grep "xsistence")
	if [ -n "$persistent_exists" ]; then
		sleep 5
		sync
		# umount every possible mounted filesystem (/home and /mnt/xsistence)
		mnt_xsistence_img_mounted=$(cat /etc/mtab | grep "/mnt/xsistence/sabayon-xsistence.ext3")
		if [ -n "$mnt_xsistence_img_mounted" ]; then
			mozilla_cache=$(find /home/*/.mozilla -nowarn -name Cache)
			for word in $mozilla_cache; do
				rm -rf $word
			done
			sync
			sleep 1
			umount -l -r -d /mnt/xsistence/sabayon-xsistence.ext3
		fi
		mnt_xsistence_mounted=$(cat /etc/mtab | grep "/mnt/xsistence")
		if [ -n "$mnt_xsistence_mounted" ]; then
			umount -l -r /mnt/xsistence
		fi
		# redundant...
		home_mounted=$(cat /etc/mtab | grep "/home ")
		if [ -n "$home_mounted" ]; then
			umount -l -r -d /home
		fi

	fi

}

start() {

	my_setup_desktop
	my_setup_password
	my_setup_xsistence
	# my_setup_samba
	my_setup_keymap
	my_wait_opengl_config
	# MOVED HERE TO AVOID RACE CONDITIONS ON WRITING /etc/profile.env variables
	# Check if LANG is forced by cmdline
	my_setup_locale
	# setup autologin for all the supported Login managers
	sabayon_setup_autologin
	# override autostart stuff if required:
	my_setup_session
	sabayon_setup_motd
	sabayon_setup_vt_autologin
	sabayon_setup_oem_livecd

}

stop() {
	my_setup_xsistence_stop
}

