# app-misc/sabayon-live/files/2.4/cdeject
#!/sbin/runscript

description="Eject LiveCD/DVD at live system shutdown/reboot"

depend() {
	after *
}

start() {
	/usr/libexec/cdeject.sh
}
