# app-misc/sabayon-live/files/2.4/installer-text.service
[Unit]
Description=Sabayon text installer setup
After=sabayonlive.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/installer-text.sh

[Install]
WantedBy=multi-user.target
