# app-misc/sabayonlive-tools/files/2.1.1/installer-gui
#!/sbin/runscript
# Copyright 2006 Sabayon Linux
# Distributed under the terms of the GNU General Public License v2

depend() {
    before xdm      
    after sabayonlive
}

start() {

        cmdline_installer_exist=$(cat /proc/cmdline | grep "installer-gui")
        if [ -z "$cmdline_installer_exist" ]; then
		ebegin "Sabayon Linux GUI Installer service not enabled"
                eend 0
        else
                ebegin "Sabayon Linux GUI Installer service"
		source /sbin/sabayon-functions.sh

		# Configure the GUI Installation
		
		# Enable autologin
		sabayon_setup_autologin

		# Configure Fluxbox
        	echo "[Desktop]" > /home/sabayonuser/.dmrc
        	echo "Session=fluxbox" >> /home/sabayonuser/.dmrc
	        chown sabayonuser /home/sabayonuser/.dmrc
                echo >> /home/sabayonuser/.fluxbox/init
                echo "session.screen0.rootCommand: xhost + & display -backdrop -window root /usr/share/backgrounds/sabayonlinux.png & installer" >> /home/sabayonuser/.fluxbox/init

		eend 0
	fi



}
