# app-misc/sabayonlive-tools/files/2.1.1/installer-text
#!/sbin/runscript
# Copyright 2006 Sabayon Linux
# Distributed under the terms of the GNU General Public License v2

depend() {
    after *      
}

start() {

        cmdline_installer_exist=$(cat /proc/cmdline | grep "installer-text")
        if [ -z "$cmdline_installer_exist" ]; then
		ebegin "Sabayon Linux Text Installer service not enabled"
                eend 0
        else
                ebegin "Sabayon Linux Text Installer service"
	        eend 0

		# switch to verbose mode
                splash_manager -c set -t default -m v &> /dev/null
		reset
		chvt 1
		clear

                echo "Welcome to Sabayon Linux Text installation." >> /etc/motd
                echo "root password: root" >> /etc/motd
                echo "to run the installation type: installer <and PRESS ENTER>" >> /etc/motd

		
	fi



}
