# net-misc/xrdp/files/xrdp-initd
#!/sbin/runscript
# Copyright 1999-2006 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/xrdp/files/xrdp-initd,v 1.1 2006/08/02 00:58:36 chutzpah Exp $

depend() {
	use net logger
}

start() {
	cd /usr/LIBDIR/xrdp
	ebegin "Starting X remote desktop daemon"
	start-stop-daemon --start --pidfile /var/run/xrdp.pid --exec /usr/LIBDIR/xrdp/xrdp -- >> ${XRDP_LOGFILE}
	ret=${?}
	eend ${ret}
	[[ ${ret} != 0 ]] && return ${ret}

	ebegin "Starting remote desktop session manager"
	start-stop-daemon --start --pidfile /var/run/sesman.pid --exec /usr/LIBDIR/xrdp/sesman >> ${SESMAN_LOGFILE}
	ret=${?}
	eend ${ret}
	return ${ret}
}

stop() {
	ebegin "Stopping X remote desktop daemon"
	start-stop-daemon --stop --pidfile /var/run/xrdp.pid
	ret=${?}
	eend ${ret}

	ebegin "Stopping remote desktop session manager"
	start-stop-daemon --stop --pidfile /var/run/sesman.pid
	ret=${?}
	eend ${ret}
	return ${ret}
}
