# net-misc/openswan/files/ipsec-initd
#!/sbin/runscript
# Copyright 1999-2008 Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2
# $Header: /var/cvsroot/gentoo-x86/net-misc/openswan/files/ipsec-initd,v 1.1 2008/09/21 12:42:32 mrness Exp $

depend() {
	need net logger
	after dns
}

start() {
	ebegin "Starting IPSec"
	ipsec setup --start
	eend $?
}

stop() {
	ebegin "Stopping IPSec"
	ipsec setup --stop
	eend $?
}

status() {
	if [ "$RC_QUIET" != "yes" ]; then
		ipsec setup --status
	fi
}
