# media-tv/mythtv/files/mythfilldatabase.cron
#!/bin/sh

QTDIR=/usr/qt/3
/usr/bin/mythfilldatabase --no-delete --quiet
