# app-misc/sabayon-mce/files/1.1/init.d/sabayon-mce
#!/sbin/runscript

depend() {
	after mtab sabayonlive
	before xdm
	need dbus
	need net
}

start() {
	ebegin "Configuring Sabayon Media Center if needed"
	/usr/libexec/sabayon-mce.sh
	eend $?
}

