# app-misc/sabayon-mce/files/1.1/init.d/sabayon-mce
#!/sbin/runscript
# Copyright 2012 Sabayon

depend() {
	# In this way autologin settings are overwritten
	after mtab sabayonlive
	before xdm
	need net
}

start() {

	SABAYON_USER="sabayonmce"
	source /sbin/sabayon-functions.sh

	if sabayon_is_mce; then

		ebegin "Sabayon Media Center mode enabled"

		# setup default .dmrc
		echo "[Desktop]" > /var/sabayonmce/.dmrc
		echo "Session=sabayon-mce" >> /var/sabayonmce/.dmrc
		chown sabayonmce /var/sabayonmce/.dmrc
		if [ -x "/usr/libexec/gdm-set-default-session" ]; then
			# oh my fucking glorious god, this
			# is AccountsService bullshit
			# cross fingers
			/usr/libexec/gdm-set-default-session sabayon-mce
		fi

		sabayon_setup_autologin

	elif ! sabayon_is_live && ! sabayon_is_mce; then
		ebegin "Sabayon Media Center mode disabled"
		sabayon_disable_autologin
	fi

	eend 0

}

