# app-misc/sabayonlive-tools/files/2.4/installer-text.service
[Unit]
Description=Sabayon text installer setup
Before=systemd-user-sessions.service

[Service]
Type=oneshot
RemainAfterExit=true
ExecStart=/usr/libexec/installer-text.sh

[Install]
WantedBy=multi-user.target