# app-misc/sabayonlive-tools/files/2.4/sabayonlive
#!/sbin/runscript
# Copyright 2004-2013 Sabayon
# Distributed under the terms of the GNU General Public License v2

depend() {
	after x-setup
	before xdm
}

start() {
	ebegin "Preparing Live system..."
	/usr/libexec/sabayonlive.sh
	eend $?
}
