# sys-libs/glibc/files/nscd.socket
[Unit]
Description=Name Service Cache Daemon Socket

[Socket]
ListenDatagram=/var/run/nscd/socket

[Install]
WantedBy=sockets.target
