# net-misc/hylafax/files/hylafax.cron
0 * * * *	root	/usr/sbin/faxqclean
25 23 * * *	root	sh /usr/sbin/faxcron | mail FaxMaster
